# etc/profile.d/my-preference.sh
# /etc/profile.d/my-preference.conf

# Colorful ls
alias ls='ls --color=auto'

# Use Vim!
EDITOR="vim"
export EDITOR
