# etc/profile.d/my-preference.sh
# /etc/profile.d/my-preference.conf

# Use Vim!
EDITOR="vim"
export EDITOR
